# test/test-configurations/shinken/etc/supervisor/conf.d/shinken-arbiter.conf
[program:shinken-arbiter]
command=/usr/bin/shinken-arbiter -c %(ENV_HOME)s/mod-webui/test/test-configurations/shinken/etc/shinken.cfg
autostart=true
user=shinken
log_stdout=false
log_stderr=true
logfile=%(ENV_HOME)s/mod-webui/test/test-configurations/shinken/log/supervisor_shinken_arbiter.log
